# gyeongju-web-api.service
[Unit]
Description=Gyeongju Travel Web API (Chat with Poser)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/gyeongju-bot/src
ExecStart=/home/ubuntu/gyeongju-bot/venv/bin/python3 web_api.py
Restart=on-failure
RestartSec=10
EnvironmentFile=/home/ubuntu/gyeongju-bot/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
