# gyeongju-bot.service
[Unit]
Description=Gyeongju Travel Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/gyeongju-bot
ExecStart=/usr/bin/python3 src/bot.py
Restart=on-failure
RestartSec=10
EnvironmentFile=/home/ubuntu/gyeongju-bot/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
